# startup/systemd/keeb_os_probe.service
[Unit]
Description=Run host OS keyboard probe
[Service]
Type=oneshot
ExecStart=/usr/local/bin/keeb_os_probe
[Install]
WantedBy=multi-user.target
